# obinit/tests/functional/res/layers/first-test-layer.obld/root/etc/profile.d/99-20190202-1-first-test-layer.sh
echo -e "\n\n-----\nExtra hello from the FIRST test layer\n-----\n\n"
