# obinit/tests/functional/res/layers/second-test-layer.obld/root/etc/profile.d/99-20190202-1-second-test-layer.sh
echo -e "\n\n-----\nExtra hello from the SECOND test layer\n-----\n\n"
