# obinit/tests/functional/res/layers/alt-root.obld/root/etc/profile.d/98-greeting.sh
echo -e "\n\n-----\nCommon hello from the THIRD test layer\n-----\n\n"
