# obinit/tests/functional/res/layers/alt-root.obld/root/etc/profile.d/99-20190330-2-third-test-layer.sh
echo -e "\n\n-----\nExtra hello from the THIRD test layer\n-----\n\n"
